# com.nodecontroller.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.nodecontroller.daemon</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/NodeController/bin/node-controller</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Library/Logs/NodeController/error.log</string>
    <key>StandardOutPath</key>
    <string>/Library/Logs/NodeController/output.log</string>
    <key>WorkingDirectory</key>
    <string>/Applications/NodeController</string>
    <key>ProcessType</key>
    <string>Background</string>
    <key>Nice</key>
    <integer>-10</integer>
</dict>
</plist> 